# setup/crontab
*/5 * * * * root echo -n "FLAG{`cat /dev/urandom | tr -dc 'a-zA-Z0-9' | fold -w 32 | head -n 1`}" > /home/simulator/flag 2>&1

